# systemd/ufw-ssh-switch.timer
[Unit]
Description=ufw-ssh-switch

[Timer]
OnUnitActiveSec=120s
OnBootSec=60s

[Install]
WantedBy=timers.target

